# devices/asus_rt-ax89x/files.edma-nss/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# The NSS data plane is brought up by the 'nss' service (/etc/init.d/nss ->
# /usr/sbin/nss-up); its output lands in logread. To stay on the plain host
# stack: uci set nss.general.enabled='0'; uci commit nss (then reboot).

exit 0
